# docker/s6-rc.d/update/finish
#!/command/execlineb -P
echo "finish update"
