# docker/s6-rc.d/publilc/finish
#!/command/execlineb -P
echo "finish public"
